# contribution/trunk/extension/ezgroupeventcalendar/doc/examples/crontab_notification.crontab
MAILTO=root@localhost
58 0,11,23 * * * /var/www/ezpublish.dom/html/ezgroupeventcalendar/cron/cron_alerts.php >> /var/www/ezpublish.dom/html/bin/logs/ezgroupeventcalendar_notifications.cron.log;

